# templates/default/aybu_create_cgroups.upstart.erb
description "Aybu create cgroups"
start on runlevel [2345]

task

script
<% node['aybu']['cgroup_controllers'].each do |ctrl| -%>
<% path = "#{@root}/#{ctrl}/#{node['aybu']['cgroup_rel_path']}" -%>
mkdir -p <%= path -%> -m 775
chown <%= @user -%>:<%= @group -%> <%= path %>
<% end -%>
end script
